# deploy/aniflow.service
[Unit]
Description=AniFlow subscription and download manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www
Group=www
WorkingDirectory=/www/wwwroot/aniflow
Environment=ANIFLOW_DATA_DIR=/www/wwwroot/aniflow/data
Environment=ANIFLOW_DOWNLOAD_DIR=/www/wwwroot/aniflow/downloads
ExecStart=/www/wwwroot/aniflow/.venv/bin/uvicorn aniflow.app:app --host 127.0.0.1 --port 8765
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
